# scripts/mbus-bus-power.service
[Unit]
Description=Enable M-Bus HAT 36V bus power via GPIO 26
After=multi-user.target

[Service]
Type=simple
ExecStart=/bin/sh -c 'gpioset -c /dev/gpiochip0 26=1; sleep infinity'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
